# system/etc/profile.d/default-editor.sh
if [ -x "/usr/bin/micro" ] ; then
    export EDITOR="micro"
  elif [ -x "/usr/bin/nano" ] ; then
	  export EDITOR="nano"    
fi;
